# init/bitlbee@.service.in
[Unit]
Description=BitlBee Per-Connection Server
After=syslog.target

[Service]
ExecStart=@sbindir@/bitlbee -I
StandardInput=socket
User=bitlbee
