# init/bitlbee.service.in
[Unit]
Description=BitlBee IRC/IM gateway

[Service]
ExecStart=@sbindir@bitlbee -F -n
KillMode=process
StateDirectory=bitlbee
RuntimeDirectory=bitlbee
User=bitlbee

[Install]
WantedBy=multi-user.target
