# init/bitlbee.service.in
[Unit]
Description=BitlBee IRC/IM gateway
After=syslog.target

[Service]
ExecStart=@sbindir@/bitlbee -F -n

[Install]
WantedBy=multi-user.target
