# init/bitlbee@.service.in
[Unit]
Description=BitlBee Per-Connection Server

[Service]
ExecStart=@sbindir@bitlbee -I
StandardInput=socket
StandardError=syslog
User=bitlbee
StateDirectory=bitlbee
RuntimeDirectory=bitlbee
