# debian/bitlbee.prerm
#!/bin/sh -e

if [ "$1" != "upgrade" ]; then
	invoke-rc.d bitlbee stop || exit 0
fi
